# systemd/neko-xray.service
[Unit]
Description=Neko Xray (VLESS REALITY Vision and XHTTP)
After=network-online.target neko-caddy.service
Wants=network-online.target
Requires=neko-caddy.service

[Service]
Type=simple
User=neko-proxy
Group=neko-proxy
ExecStart=/usr/local/libexec/neko/xray run -c /etc/neko/config/xray.json
Restart=on-failure
RestartSec=3s
LimitNOFILE=1048576
CapabilityBoundingSet=
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
UMask=0027

[Install]
WantedBy=multi-user.target
